# osx/dbus/org.freedesktop.dbus-session.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>org.freedesktop.dbus-session</string>

	<key>ProgramArguments</key>
	<array>
		<string>@HEXCHAT_BASE_PATH@/Contents/Resources/bin/dbus-daemon</string>
		<string>--nofork</string>
		<string>--config-file=@HEXCHAT_BASE_PATH@/Contents/Resources/etc/dbus-1/session.conf</string>
	</array>

	<key>Sockets</key>
	<dict>
		<key>unix_domain_listener</key>
		<dict>
			<key>SecureSocketWithKey</key>
			<string>DBUS_LAUNCHD_SESSION_BUS_SOCKET</string>
		</dict>
	</dict>

</dict>
</plist>
